# alsa-fix/91-soundblasterx-g6.rules
# SoundBlasterX G6 - Fix ALSA defaults on card detection
# Matches on the sound subsystem (not raw USB), so the ALSA card is already registered.
# ATTRS walks up to the USB parent to match vendor/product ID.
SUBSYSTEM=="sound", ACTION=="add", DEVPATH=="*/controlC*", \
  ATTRS{idVendor}=="041e", ATTRS{idProduct}=="3256", \
  RUN+="/bin/sh -c '/usr/local/bin/soundblasterx-g6-alsa-fix.sh &'"
